# LaunchAgents/ai.hermes.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.hermes.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lyndonarthurson/.hermes/hermes-agent/venv/bin/python</string>
        <string>-m</string>
        <string>hermes_cli.main</string>
        <string>gateway</string>
        <string>run</string>
        <string>--replace</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/lyndonarthurson/.hermes/hermes-agent</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/lyndonarthurson/.hermes/hermes-agent/venv/bin:/Users/lyndonarthurson/.hermes/hermes-agent/node_modules/.bin:/opt/homebrew/Cellar/node/25.5.0/bin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/opt/homebrew/bin:/Users/lyndonarthurson/.antigravity/antigravity/bin:/opt/homebrew/opt/node@22/bin:/Applications/Docker.app/Contents/Resources/bin:/Users/lyndonarthurson/.local/bin:/opt/homebrew/sbin</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/lyndonarthurson/.hermes/hermes-agent/venv</string>
        <key>HERMES_HOME</key>
        <string>/Users/lyndonarthurson/.hermes</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/lyndonarthurson/.hermes/logs/gateway.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/lyndonarthurson/.hermes/logs/gateway.error.log</string>
</dict>
</plist>
